# Library/LaunchAgents/my.sync.ocr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>my.screenshot.ocr</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>exec "$HOME/Applications/Screenshots Sync.app/Contents/MacOS/Automator Application Stub"</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/alex/Screenshots/Processing</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/Users/alex/Library/Logs/screenshots-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alex/Library/Logs/screenshots-sync.log</string>
</dict>
</plist>
